# suse-online-update.timer
[Unit]
Description=Daily non-interactive online update

[Timer]
OnCalendar=daily
AccuracySec=24h
Persistent=true

[Install]
WantedBy=timers.target
